# scripts/websocket-proxy-setup.service
[Unit]
Description=Setup WebSocket Proxy for Nginx
After=network.target nginx.service
Before=nginx.service

[Service]
Type=oneshot
ExecStart=/bin/bash /root/vpn/setup-websocket-proxy.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
